# files/etc/supervisor/conf.d/apache2.ini
[program:apache2]
command = /usr/sbin/apache2ctl -DFOREGROUND
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stdout_logfile = /var/log/supervisor/apache2-stdio.log
stdout_logfile_maxbytes = 4MB
redirect_stderr = true
