# files/etc/supervisor/conf.d/cupsd-helper.ini
[program:cupsd-helper]
command = /usr/local/sbin/cups-helper listener
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/cupsd-helper.log
stdout_logfile_maxbytes = 4MB
