# files/etc/supervisor/conf.d/cupsd.ini
[program:cups]
command = /usr/sbin/cupsd -l
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/cupsd.log
stdout_logfile_maxbytes = 4MB
