# files/etc/supervisor/conf.d/dbus.ini
[program:dbus]
command = /usr/local/sbin/dbus-runner
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/dbus.log
stdout_logfile_maxbytes = 4MB
